# surface_hardware/rootfs/etc/services.d/surface-hardware/run
#!/usr/bin/with-contenv bashio

exec python3 /usr/bin/hardware.py